# services/atp-todo-recur.service
[Unit]
Description=ATP Todo Recur - Daily recurring todo processing
After=network.target

[Service]
Type=oneshot
Environment="ATP_DIR=${ATP_DIR}"
Environment="PATH=${SCRIPTS}:${GOPATH}/bin:${HOME}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin"
ExecStart=${SCRIPTS}/atp todo recur