# services/watch-knowledge.service
[Unit]
Description=Watches the knowledge repo for changes and commits them

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/home/arjun/repos/github.com/arjungandhi/dot/scripts/autocommit /home/arjun/repos/github.com/arjungandhi/knowledge

[Install]
WantedBy=default.target
